# roles/os_base/files/apt-daily.timer.j2
[Timer]
OnCalendar=
OnCalendar={{ debian_automatic_update_download_schedule }}
RandomizedDelaySec=5m
Persistent=true
