# roles/os_base/files/apt-daily-upgrade.timer.j2
[Timer]
OnCalendar=
OnCalendar={{ automatic_update_schedule }}
RandomizedDelaySec=5m
Persistent=true
