# roles/os_base/files/apt-daily.timer.j2
[Unit]
Description=Daily apt download activities

[Timer]
OnCalendar=*-*-* {{ automatic_update_download_time }}
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target
