# roles/os_base/files/auter.timer.j2
[Unit]
Description=Automatically runs auter every day
Requires=auter.service

[Timer]
Unit=auter.service
OnCalendar=*-*-* {{ automatic_update_download_time }}

[Install]
WantedBy=timers.target
